# librespeed-mqtt.service
[Unit]
Description=LibreSpeed WAN test published to MQTT
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/librespeed-mqtt --config /etc/librespeed-mqtt/config.toml
